# cookbooks/mediacenter/templates/ubuntu/sonarr.service.erb
[Unit]
Description=Sonarr Daemon
After=syslog.target network.target

[Service]
User=<%= @user %>
Group=<%= @group %>

Type=simple
ExecStart=/usr/bin/mono /opt/NzbDrone/NzbDrone.exe -nobrowser -data=/home/kodi/.config/sonarr
TimeoutStopSec=20
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
